# docker/supervisord.fly.conf
# supervisord.fly.conf — Fly.io multi-process supervision for Hermes.
#
# Fly runs the container entrypoint as a CHILD of Fly's own init, so the
# image's s6-overlay /init can't be PID 1 and crash-loops. We bypass s6 (see
# fly.toml [experimental] entrypoint) and use supervisord instead to run the
# two long-lived processes a remote-desktop-attachable instance needs:
#
#   gateway    — `hermes gateway run` (Telegram long-poll + kanban dispatcher)
#   dashboard  — `hermes dashboard`  (HTTP control plane the desktop app attaches to)
#
# Fly's init still reaps zombies and forwards signals (PID-1 duties); supervisord
# adds the one thing Fly's init does NOT do — per-process restart-on-crash.
#
# supervisord runs as root (started by the entrypoint after it chowns the
# volume) and drops each child to the hermes user (uid 10000) via `user=`.
# HOME=/opt/data is REQUIRED per child — without it, uid 10000 writes lock/state
# under /root and EACCESes. Children inherit Fly secrets from supervisord's env.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/opt/data/supervisord.pid
# supervisord itself stays root so it can setuid children to hermes.
user=root

[program:gateway]
command=/opt/hermes/.venv/bin/hermes gateway run
directory=/opt/data
user=hermes
environment=HOME="/opt/data"
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=25
stopsignal=TERM
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:dashboard]
command=/opt/hermes/.venv/bin/hermes dashboard --no-open --skip-build --host 0.0.0.0 --port 9119
directory=/opt/data
user=hermes
environment=HOME="/opt/data"
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
